# ansible/roles/system/rate_mirrors/templates/rate-mirrors.timer.j2
[Unit]
Description=Refresh pacman mirrorlist with rate-mirrors

[Timer]
OnCalendar={{ rate_mirrors_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
