# ansible/roles/system/restic/templates/restic-r2.timer.j2
[Unit]
Description=Restic backup to Cloudflare R2

[Timer]
OnCalendar={{ restic_backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
